# etc/profile.d/spring-cli.sh
#!/bin/bash
export SPRING_HOME=/opt/spring-1.2.5.RELEASE
PATH=${SPRING_HOME}/bin:${PATH}
export PATH
#. /opt/spring-1.2.5.RELEASE/shell-completion/bash/spring
