# etc/profile.d/user-bin.sh
# add ~/bin to path for all users.
PATH="${HOME}/bin:${PATH}"
export PATH

